# files/root/turn_on_accesspoint.sh
#!/bin/sh /etc/rc.common

uci set wireless.default_radio0.disabled=0
uci commit wireless
wifi  up

